# cnetflow.service
[Unit]
Description=Cnetflow service.

[Service]
ExecStart=/usr/local/cnetflow/cnetflow

[Install]
WantedBy=multi-user.target