# build/_CPack_Packages/Linux/DEB/cnetflow-0.1.1-Linux/etc/systemd/system/cnetflow.service.d/local.conf
[Service]
Environment="PG_CONN_STRING=postgresql://postgres.your-tenant-id:your-super-secret-and-long-postgres-password@192.168.100.78:6543/postgres"
Environment="UV_THREADPOOL_SIZE=16"
